# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# Auto mount disk drive on startup
# change to your drive path and uncomment below
#mount_hdd /dev/sda1 /mnt/sda1

exit 0
